# systemd/chowdyd.socket
[Unit]
Description=chowdy authentication and management sockets

[Socket]
# Slot 0 — auth socket (PAM module). Daemon reads SO_PEERCRED and only
# accepts requests from uid 0 (PAM in a root-running login service).
ListenStream=/run/chowdy/auth.sock
SocketUser=root
SocketGroup=chowdy
SocketMode=0660

# Slot 1 — management socket (CLI). Daemon uses SO_PEERCRED to scope every
# request to the caller's uid.
ListenStream=/run/chowdy/mgmt.sock
SocketUser=chowdy
SocketGroup=chowdy
SocketMode=0666

Accept=false

[Install]
WantedBy=sockets.target
